# docker/php-fpm/crontab
0 8 * * * php /app/bin/console app:uploadNewDataFromCfpWikiCommand >> /var/log/schedule.log 2>&1